# packaging/ubuntu-trusty/debian/amazon-ecs-init.postinst
if [ -f /bin/systemctl ]; then
    /bin/systemctl daemon-reload
    /bin/systemctl enable ecs
    /bin/systemctl start ecs
fi
